# {{cookiecutter.repo_name}}/provision/conf/gunicorn_init.conf
# log: sudo vim var/log/upstart/gunicorn_init.log 

description "Gunicorn started with upstart"

start on (filesystem)
stop on runlevel [016]

console log

chdir /webapps/{{cookiecutter.repo_name}}/{{cookiecutter.project_name}}
exec /webapps/.virtualenvs/{{cookiecutter.project_name}}/bin/gunicorn {{cookiecutter.project_name}}.wsgi:application  
